# linux-basic/linux-basic/crontab/lv.script
#!/bin/bash

TARGET="/mnt/first-lv-data"
LOG="$TARGET/copy.log"

mkdir -p "$TARGET/usr_sample"

rsync -av /usr/share/doc/ "$TARGET/usr_sample/" >> "$LOG" 2>&1

echo "$(date): copy completed" >> "$LOG"

